# app-emulation/anbox/files/anbox-init.conf
# Copyright (c) 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "set alsa default config"
author        "tean@fydeos.io"

start on started system-services

script
  #anbox-container --daemon --privileged --data-path=/mnt/stateful_partition/var_overlay/lib/anbox
  anbox-container /mnt/stateful_partition/var_overlay/lib/anbox
end script